# samba/supervisord.conf
[supervisord]
nodaemon=true

[program:smbd]
command=/usr/sbin/smbd -FS
autorestart=true

[program:nmdb]
command=/usr/sbin/nmbd -FS
autorestart=true

